# pkg/deb/DEBIAN/prerm
#!/bin/sh -e
systemctl disable --now peerbanhelper.service 2> /dev/null
